# deploy/linux/laplace-api.service
[Unit]
Description=Laplace API (OpenAI-compat endpoints + self-hosted SPA)
Documentation=https://github.com/SaltyPatron/Laplace
After=network.target laplace-postgresql.service
Wants=laplace-postgresql.service

[Service]
Type=simple
User=laplace-runner
Group=laplace-runner
WorkingDirectory=/opt/laplace/app
# Non-secret runtime config (DB, perfcache, URLs). Owned by laplace-runner; CI may refresh.
EnvironmentFile=/opt/laplace/app/laplace-api.env
# Operator-published secrets (setup-host seeds; CI publish may refresh). Optional (-).
EnvironmentFile=-/opt/laplace/secrets/lichess.env
EnvironmentFile=-/opt/laplace/secrets/stripe.env
ExecStart=/usr/bin/dotnet /opt/laplace/app/Laplace.Endpoints.OpenAICompat.dll
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
SyslogIdentifier=laplace-api

[Install]
WantedBy=multi-user.target
